# systemd/.config/systemd/user/check-battery.service
[Unit]
Description=Notify on low battery

[Service]
ExecStart=/home/lervag/scripts/bin/check-battery.sh
